# scripts/mytasks.cron
#! /bin/bash

/var/www/mytasks/mytasks backup
